# root/etc/supervisor.d/traefik.ini
[program:traefik]
command=/opt/traefik/bin/traefik --configFile=/opt/traefik/etc/traefik.toml